# etc/supervisor/mosn.conf
[program:mosn]
command=mosnd start -c /home/admin/mosn/conf/mosn.conf
user=admin
startsecs=1
autostart=true
autorestart=true
stdout_events_enabled = true
stderr_events_enabled = true