# etc/supervisor/mosn.conf
[program:mosn]
command=mosnd -c /home/admin/mosn/conf/mosn.conf
user=root
autostart=false
autorestart=false
stdout_events_enabled = true
stderr_events_enabled = true